# dot_config/autostart/kopia-ui.desktop
[Desktop Entry]
Type=Application
Exec=/opt/KopiaUI/kopia-ui
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
Name=KopiaUI

